# deploy/pain-signals.service
[Unit]
Description=Pain Signals - Reddit Problem Discovery
After=network.target

[Service]
User=kirill
Group=kirill
WorkingDirectory=/var/www/pain-signals
ExecStart=/var/www/pain-signals/venv/bin/gunicorn --bind 127.0.0.1:5001 --workers 2 --timeout 120 app:app
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
